# supervisor/bluetooth_service.conf
[program:bluetooth_service]
command=/usr/bin/python3 ble_peripheral.py
directory=/home/pi/connected-devices/Lampi/bluetooth/
priority=500
autostart=true
autorestart=true