# supervisor/lampi_app.conf
[program:lampi_app]
command=/usr/bin/python3 lampi_app.py
directory=/home/pi/lampi-looper/Lampi/
priority=400
user=pi
autostart=true
autorestart=true