# supervisor/lampi_db.conf
[program:lampi_db]
command=/usr/bin/python3 lampi_db.py
directory=/home/pi/lampi-looper/Lampi/
priority=200
user=root
autostart=true
autorestart=true