# supervisor/pigpiod.conf
[program:pigpiod]
command=/usr/local/bin/pigpiod -l -g
priority=100
autostart=true
autorestart=true