# configs/bluetooth_service.conf
[program:bluetooth_service]
command=/connected-devices/Lampi/bluetooth/ble_peripheral.py
directory=/home/pi
priority=500
autostart=true
autorestart=true