# configs/lampi_ui.conf
[program:lami_ui]
command=/usr/bin/python3 lampi_ui.py
directory=/home/pi/connected-devices/Lampi/
priority=400
user=pi
autostart=true
autorestart=true