# configs/lampi_ui.conf
[program:lami_ui]
command=lampi_ui.py
directory=
priority=400
user=pi
autostart=true
autorestart=true